# nos/datapath/rootfs-overlay/etc/systemd/system/edgenos-l3.service
[Unit]
Description=EdgeNOS-4610 L3 config (addrs + static routes from /etc/edged)
After=bcmd.service
Wants=bcmd.service
ConditionPathExists=/etc/edged/addrs.conf

[Service]
Type=oneshot
ExecStart=/usr/sbin/edgenos-l3-config.sh
RemainAfterExit=yes
TimeoutStartSec=180

[Install]
WantedBy=multi-user.target
